# jobs/awslogs/templates/awslogs_log_rotate.cron.erb
# Version: 1.1.2-rpm
MAILTO=""
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

30 * * * * root logrotate -s /var/log/logstatus /etc/logrotate.d/awslogs
